# config/chroot_local-includes/lib/systemd/system/gobby-server.service
[Unit]
Description=Gobby Server Daemon
After=network.target

[Service]
User=infinoted
Group=infinoted
ExecStart=/usr/bin/infinoted
